# configs/supervisor.ini
[program:python-alarm-realtime]
directory=/data/www/slack-aliyun/configs
command=/usr/bin/python /data/www/slack-aliyun/alarm-realtime.py
autostart=true
autorestart=true
stopsignal=QUIT
killasgroup=true
stdout_logfile=/data/logs/slack-aliyun/alarm-realtime.stdout.log
stderr_logfile=/data/logs/slack-aliyun/alarm-realtime.stderr.log